# raspberry-blue-bridge.service

[Unit]
Description=Raspberry Blue Bridge

[Service]
ExecStart=_BIN_/raspberry-blue-bridge

[Install]
WantedBy=multi-user.target
